# openwrt/luci-app-ninja/root/etc/uci-defaults/luci-opengpt
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ninja[-1]
	add ucitrack ninja
	set ucitrack.@ninja[-1].init=ninja
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
